# src/tracker-control/tracker-control.service.in
[Unit]
Description=Tracker Filesystem Indexing Proxy
ConditionUser=!root

[Service]
Type=dbus
BusName=org.freedesktop.Tracker3.Miner.Files.Control
ExecStart=@libexecdir@/tracker-miner-fs-control-3
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
